# ulogd/systemd/ulog.socket
[Unit]
Requires=ulogd.service

[Socket]
ListenFIFO=/dev/ulog
#ListenStream=/dev/ulog

[Install]
WantedBy=sockets.target

