# etc/systemd/servicewall-logs.socket
[Socket]
ListenFIFO=/dev/ulog
SocketMode=0600
#ListenStream=/dev/ulog

[Install]
WantedBy=sockets.target

